# tailscale/tailscale.container
[Unit]
Description=Tailscale Podman Container
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Container]
Image=docker.io/tailscale/tailscale:latest
ContainerName=tailscale
AutoUpdate=registry
Network=host
DropCapability=ALL
AddCapability=NET_ADMIN
AddCapability=NET_RAW
AddCapability=SYS_MODULE
PodmanArgs=--security-opt=no-new-privileges
PodmanArgs=--memory=512m
PodmanArgs=--pids-limit=256
ReadOnly=true
Environment=TS_HOSTNAME=%H
Environment=TS_STATE_DIR=/var/lib/tailscale
Volume=/dev/net/tun:/dev/net/tun
Volume=/var/lib/tailscale:/var/lib/tailscale:Z
Volume=tmpfs:/tmp
Volume=tmpfs:/run
HealthCmd=/usr/local/bin/tailscale status --json
HealthInterval=60s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

[Install]
WantedBy=default.target
